# deploy/systemd/pcloud-migration-supervisor.timer
[Unit]
Description=Advance the pCloud to Proton migration every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=pcloud-migration-supervisor.service

[Install]
WantedBy=timers.target
